# deploy/xboard.service
[Unit]
Description=XBoard Panel Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xboard/panel
EnvironmentFile=-/etc/default/xboard
ExecStart=/opt/xboard/panel/xboard serve --config /opt/xboard/panel/config.yml
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
